# crontab
*/10 * * * * umask 0022; cd ~/fetch-win32-symbols; python start.py
*/5 * * * * umask 0022; cd ~/socorro; PYTHONPATH=.:thirdparty ./socorro-virtualenv/bin/python socorro/cron/crontabber.py
